# s6/php-fpm/run
#!/bin/sh

set -e

umask 0077

mkdir -p \
  /data/torrents \
  /data/users \
  /tmp/php \
  /tmp/php/session

/usr/sbin/php-fpm8.2 --nodaemonize --fpm-config /config/php-fpm/php-fpm.conf
