# s6/rtorrent/run
#!/bin/sh

mkdir -p /tmp/rtorrent
rm -f /data/session/rtorrent.lock

exec /usr/bin/rtorrent -n -o import=/config/rtorrent/rtorrent.rc
